# vision/waffles_vision.service
[Service]
Type=simple
User=waffle
ExecStart=waffles_vision

[Install]
WantedBy=multi-user.target
